# services/appSave.timer
[Unit]
Description=Timer to update Database

[Timer]
OnBootSec=1s
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
